# files/apt-daily-upgrade.timer.d/fk.override.conf
[Timer]
OnCalendar=
RandomizedDelaySec=25
OnBootSec=5min
