# files/apt-daily-upgrade.timer.d/ziska.override.conf
[Timer]
OnStartupSec=60s
RandomizedDelaySec=5s