# files/apt-daily.timer.d/ziska.override.conf
[Timer]
OnStartupSec=25s
RandomizedDelaySec=5s